# crontab
*/60 * * * * vdirsyncer discover > /dev/null
*/3 * * * * vdirsyncer sync > /dev/null
